# email-backup/email-backup-cleanup.service
[Unit]
Description=Email Backup Cleanup Service
After=network-online.target docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/home/teo-mcarthur/home-server/email-backup
ExecStart=/usr/bin/docker compose run --rm email-backup cleanup.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
